# templates/labkey_server.service.j2
[Unit]
Description=LabKey Server Application
After=syslog.target network.target

[Service]
Type=simple
Environment="JAVA_HOME=/usr/local/java"
Environment="LABKEY_HOME={{ labkey_home }}"
WorkingDirectory={{ labkey_home }}
OOMScoreAdjust=-500

{# Výpočet heap: pokud není ručně nastaveno, počítá se z RAM.
   RAM > 8 GB → max = 75 % RAM, jinak max = 50 % RAM.
   Min heap = polovina max heap. #}
{% set _ram_mb = ansible_memtotal_mb %}
{% if labkey_jvm_max_heap != "" %}
{% set _max_heap = labkey_jvm_max_heap %}
{% elif _ram_mb > 8192 %}
{% set _max_heap = ((_ram_mb * 0.75) | int | string) + "M" %}
{% else %}
{% set _max_heap = ((_ram_mb * 0.50) | int | string) + "M" %}
{% endif %}
{% if labkey_jvm_min_heap != "" %}
{% set _min_heap = labkey_jvm_min_heap %}
{% else %}
{% set _min_heap = ((_max_heap | regex_replace('[^0-9]', '') | int // 2) | string) + "M" %}
{% endif %}
ExecStart=/usr/local/java/bin/java \
    -Xms{{ _min_heap }} \
    -Xmx{{ _max_heap }} \
    -XX:+HeapDumpOnOutOfMemoryError \
    -Djava.io.tmpdir={{ labkey_home }}/labkey-tmp \
    -Dorg.apache.catalina.startup.EXIT_ON_INIT_FAILURE=true \
    -DterminateOnStartupFailure=true \
    --add-opens=java.base/java.lang=ALL-UNNAMED \
    --add-opens=java.base/java.io=ALL-UNNAMED \
    --add-opens=java.base/java.util=ALL-UNNAMED \
    --add-opens=java.base/java.util.concurrent=ALL-UNNAMED \
    --add-opens=java.rmi/sun.rmi.transport=ALL-UNNAMED \
    -jar labkeyServer.jar

SuccessExitStatus=0 143
Restart=on-failure
RestartSec=15

User={{ labkey_system_user }}
Group={{ labkey_system_user }}

[Install]
WantedBy=multi-user.target
